# elasticsearch/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/usr/share/elasticsearch/bin/es-docker

[program:nginx]
command=/usr/sbin/nginx
nodaemon=true

